# packages/release/base/shadow/files/shadow.cron
#!/bin/sh
#
# /etc/cron.d/daily/shadow
#

/usr/sbin/grpconv
/usr/sbin/pwconv
